# faustweb.conf
# faustweb is a remote compilation web service for faust

description "Remote compilation web service for faust"
author      "Yann Orlarey"


env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/opt/android/android-ndk-r8e:/opt/android/sdk/platform-tools:/opt/android/sdk/tools
chdir /home/faust/yann/faustweb
setuid faust

start on runlevel [2345]
stop  on runlevel [016]

exec /home/faust/yann/faustweb/faustweb